# deploy/linux/monidroid.service
[Unit]
Decription=Monidroid server

[Service]
ExecStart=/opt/monidroid/monidroid-server
RestartKillSignal=SIGINT
KillSignal=SIGINT
