# deploy/linux/monidroid.service
[Unit]
Decription=Monidroid server

[Service]
ExecStart=/usr/bin/monidroid-server
RestartKillSignal=SIGINT
KillSignal=SIGINT
